# ArchPackages/linuxconsole-five/linuxconsole-five.install
pre_install() {
	systemctl stop inputattach
	systemctl disable inputattach
}

post_install() {
	systemctl start inputattach
	systemctl enable inputattach
	systemctl daemon-reload
}

pre_upgrade() {
	pre_install
}

post_upgrade() {
	post_install
}

post_remove() {
	pre_install
	systemctl daemon-reload
}
